# src/debian/prerm
#!/bin/bash

# stop the shodan service before removing the package
systemctl stop shodan.service
systemctl disable shodan.service
